# ReMailer.service
[Unit]
Description=Reformat mails from XING & Co
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/ReMailer
#ExecStartPre=/usr/local/bin/pipenv update
#TimeoutStartSec=30
ExecStart=/usr/local/bin/pipenv run python main.py

[Install]
WantedBy=multi-user.target
